# udevadm-monitor.service
[Unit]
Description=Monitor udev events

[Service]
Type=simple
ExecStart=udevadm monitor
